# board/rescue/fsoverlay/etc/profile.d/30-welcome.sh
echo
echo batocera-info 2>/dev/null
echo "OS version: $(batocera-version)"
echo
